# examples/fetch-runner.service
[Unit]
Description=fetch-runner: poll git branches and run scripts on new commits
Documentation=https://github.com/BYU-ODH/fetch-runner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# ===========================================================================
# CUSTOMIZE: edit the lines in this block for each deployment.
# ===========================================================================

# MUST match the `user` field under [general] in the jobs.toml this unit
# starts. This is the *polling* user — fetch-runner itself runs as this
# user, but it does NOT own the git repos: each repo is owned by its
# job's [[jobs]].run_as user, and `git fetch` / `git checkout` are
# dispatched as that user via sudo. The polling user only needs search
# permission on the repo paths.
# fetch-runner exits at startup if User= and [general].user disagree.
User=fetch-runner
Group=fetch-runner

# Path to the fetch-runner binary and to the config file it should load.
# If you installed fetch-runner from a venv, point ExecStart at that venv's
# bin/fetch-runner. The config path is positional.
ExecStart=/usr/local/bin/fetch-runner /etc/fetch-runner/jobs.toml

# Every directory that any process under this unit needs to write to must
# appear here — including the git repos (`git fetch` and `git checkout`
# run as the run_as user via sudo, but they are still children of this
# unit and subject to its filesystem sandbox) and any app state the
# deploy scripts touch. Separate with spaces. The rest of the filesystem
# is read-only thanks to ProtectSystem.
ReadWritePaths=/srv

# ===========================================================================
# DO NOT MODIFY below this line casually.
# These settings are the sandbox that keeps a long-running deploy service from
# having broad access to the host. If you must relax one, document why.
# ===========================================================================

Restart=on-failure
RestartSec=10s
TimeoutStopSec=30s

# Send both fetch-runner logs and deploy-script output to journald so there is
# one place to inspect failures.
StandardOutput=journal
StandardError=journal

# NoNewPrivileges and RestrictSUIDSGID are deliberately NOT set here:
# fetch-runner uses `sudo -u <run_as>` to dispatch deploy scripts to per-job
# users, and both settings block sudo's setuid transition. The narrow
# sudoers fragment (see `fetch-runner --print-sudoers`) is what bounds what
# fetch-runner can actually do with that privilege — it may only run the
# exact script paths listed, as the exact run_as users listed, with no
# password and a fixed env_keep allowlist. If every job in your jobs.toml
# uses run_as = [general].user (the single-user mode), you can re-enable:
#   NoNewPrivileges=true
#   RestrictSUIDSGID=true
CapabilityBoundingSet=
AmbientCapabilities=

# Keep the filesystem read-only by default, then punch narrow write holes back
# in with `ReadWritePaths=` above for the repos and app state that deployments
# genuinely need to modify.
ProtectSystem=strict

# Deploy scripts should not need the service user's home directory.
ProtectHome=read-only

PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true

# Prevent deploy code from reconfiguring host control groups. If your rollout
# tool genuinely needs cgroup access, revisit this deliberately.
ProtectControlGroups=true

ProtectClock=true
ProtectHostname=true

# Hide unrelated processes from the service; deploy hooks should not be
# inspecting the rest of the machine.
ProtectProc=invisible

# Kernel / namespace hardening.
LockPersonality=true
RestrictRealtime=true

# Prevent deploy code from creating new namespaces inside the sandbox.
RestrictNamespaces=true

RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Syscall filter.
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
